# build/package/systemd/cold-brew-server.service
[Unit]
Description=Cold Brew
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
ExecStart=/usr/local/bin/cold-brew-server
Restart=on-failure

[Install]
WantedBy=multi-user.target
